# services/thermald/run
#!/bin/sh
exec 2>&1
[ -x /usr/sbin/thermald ] || exit 1
exec thermald --no-daemon
